# templates/default/service.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:        iproute2
# Required-Start:    
# Required-Stop:     S 0 1 6
# Default-Start:   2
# Default-Stop:    
# Short-Description:
### END INIT INFO

# *****************************************************************
# Dynamic iproute2 configuration  --  DO NOT EDIT THIS FILE BY HAND
# Generated by chef - <%= @date %>
# *****************************************************************
PATH=/sbin:/bin:/usr/sbin:/usr/bin

. /lib/lsb/init-functions

case "$1" in
start|stop|status)
 invoke-rc.d tomcat $1 <%= @name %>
 ;;
*)
 echo "Syntaxe: $(basename $0) [start|stop|status]..." && false
esac

